# infra/goweb.service
[Unit]
Description=goweb

[Service]
Type=simple
Restart=always
RestartSec=5s
ExecStart=<DIR>/<BIN>

[Install]
WantedBy=multi-user.target